# packages/shared-state-mesh_config/files/etc/uci-defaults/shared-state_mesh_config_register
#!/bin/sh
uci set shared-state.mesh_config=dataType
uci set shared-state.mesh_config.name='mesh_config'
uci set shared-state.mesh_config.scope='community'
uci set shared-state.mesh_config.ttl='2400'
uci set shared-state.mesh_config.update_interval='30'
uci commit shared-state
